# rootfs/etc/services.d/mysqld/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh:
mysqld --user mysql
